# docker/scripts/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log

[program:Xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1920x1080x24 -dpi 96
autorestart=true
priority=1
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.log

[program:xfce4]
command=/bin/bash -c "mkdir -p /tmp/runtime-etudiant && chmod 700 /tmp/runtime-etudiant && until xdpyinfo -display :1 >/dev/null 2>&1; do sleep 0.2; done && exec dbus-launch --exit-with-session startxfce4"
user=etudiant
environment=DISPLAY=":1",HOME="/home/etudiant",USER="etudiant",XDG_RUNTIME_DIR="/tmp/runtime-etudiant",XDG_SESSION_TYPE="x11"
autorestart=true
priority=2
stdout_logfile=/var/log/xfce4.log
stderr_logfile=/var/log/xfce4.log

[program:x11vnc]
command=/bin/bash -c "until xdpyinfo -display :1 >/dev/null 2>&1; do sleep 0.2; done && exec x11vnc -display :1 -nopw -shared -forever -noxdamage -xkb -listen 0.0.0.0"
autorestart=true
priority=3
stdout_logfile=/var/log/x11vnc.log
stderr_logfile=/var/log/x11vnc.log

[program:novnc]
command=/bin/bash -c "until netstat -tlnp 2>/dev/null | grep -q ':5900'; do sleep 0.5; done && exec /usr/bin/websockify --web=/usr/share/novnc/ 6901 localhost:5900"
autorestart=true
priority=4
stdout_logfile=/var/log/novnc.log
stderr_logfile=/var/log/novnc.log

[program:disable-screensaver]
command=/bin/bash -c "sleep 20 && xset s off -dpms 2>/dev/null; true"
user=etudiant
environment=DISPLAY=":1",HOME="/home/etudiant",USER="etudiant",XDG_RUNTIME_DIR="/tmp/runtime-etudiant"
autorestart=false
priority=5
